# gentoo/binhost/service/builder/run
#!/bin/bash

if [ ! -e /portage/world ]
then
	cp -r /skel/portage/* /portage/
fi

while true
do
	inotifywait -e modify -r /portage
	sleep 1
	/usr/local/bin/builder.sh
done
